# overlay/package/utils/trousers/files/tcsd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2017 Chris Blake (chrisrblake93@gmail.com)

START=99

USE_PROCD=1
NAME=tcsd
PROG=/usr/sbin/tcsd

start_service() {
  [ -e /etc/tcsd.conf ] || return
	procd_open_instance
	procd_set_param command "$PROG"
  procd_set_param respawn
	procd_close_instance
}

stop_service() {
	service_stop "$PROG"
}

reload() {
	service_reload "$PROG"
}
